# moonshine-voice-typing.service
\
[Unit]
Description=Moonshine Voice Typing
After=graphical-session.target sound.target

[Service]
Type=simple
WorkingDirectory=%h/moonshine_voice_typing_v1
ExecStart=%h/moonshine_voice_typing_v1/.venv/bin/python -m voice_typing.daemon
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
